# hopper-bot@.service
[Unit]
Description=Hopper Discord Bot (%i)
After=network-online.target
Wants=network-online.target
ConditionPathExists=%h/%i/hopper.py
ConditionPathExists=%h/%i/.env

[Service]
Type=simple
WorkingDirectory=%h/%i
ExecStart=/usr/bin/python3 -u %h/%i/hopper.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
