# hopper-bot.service
[Unit]
Description=Hopper Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/hopper-bot
ExecStart=/usr/bin/python3 -u %h/hopper-bot/hopper.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits (optional)
# MemoryLimit=512M
# CPUQuota=50%

[Install]
WantedBy=default.target
