# hopper-backup@.timer
[Unit]
Description=Run Hopper Bot Backup every hour (%i)

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
